# infra/podman/aemr-bot-db.container
# Quadlet-юнит базы Postgres для aemr-bot (Podman 4.4+).
# Кладётся в /etc/containers/systemd/ (rootful) → генерится aemr-bot-db.service.
# Перенос свойств из infra/docker-compose.yml (сервис db) один в один.
#
# EnvironmentFile должен содержать POSTGRES_DB / POSTGRES_USER / POSTGRES_PASSWORD
# (тот же infra/.env, что и у docker-compose). Путь при необходимости поправьте.
[Unit]
Description=aemr-bot Postgres 16
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=aemr-bot-db
# Тот же pinned tag+digest, что в docker-compose.yml. Обновлять вместе со смоук-тестом.
Image=docker.io/library/postgres:16-alpine@sha256:4e6e670bb069649261c9c18031f0aded7bb249a5b6664ddec29c013a89310d50
Network=aemr-bot.network
Volume=aemr-bot-db_data:/var/lib/postgresql/data
EnvironmentFile=/home/aemr/aemr-bot/infra/.env
# pg_stat_statements должен быть в shared_preload_libraries ДО старта Postgres.
Exec=postgres -c shared_preload_libraries=pg_stat_statements -c pg_stat_statements.track=all
# Readiness-проба (для `podman healthcheck` и мониторинга). Порядок старта бота
# обеспечивается связкой After/Requires + Restart=always у бота (см. aemr-bot.container).
HealthCmd=/bin/sh -c "pg_isready -U $POSTGRES_USER"
HealthInterval=5s
HealthTimeout=3s
HealthRetries=20
# Hardening (как в compose): сбросить все capabilities, добавить только нужные
# Postgres для chown data-каталога; запретить setuid-эскалации.
DropCapability=all
AddCapability=CHOWN
AddCapability=DAC_READ_SEARCH
AddCapability=FOWNER
AddCapability=SETGID
AddCapability=SETUID
NoNewPrivileges=true
LogDriver=journald

[Service]
# unless-stopped из compose → в systemd это always (рестарт при любом падении,
# но не после ручного `systemctl stop`).
Restart=always
TimeoutStartSec=120

[Install]
WantedBy=default.target
